# docker/workspace/crontab/root
* * * * * php /var/www/karmabot/artisan schedule:run >> /dev/null 2>&1
